# deploy/locolive.service
[Unit]
Description=LocoLive Backend API
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/locolive
ExecStart=/opt/locolive/bin/api
Restart=always
RestartSec=5

# Environment file
EnvironmentFile=/opt/locolive/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=locolive

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/locolive

[Install]
WantedBy=multi-user.target
